# config/com.kb.processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kb.processor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sonseongjun/.nvm/versions/node/v20.18.1/bin/node</string>
        <string>/Users/sonseongjun/knowledge-base/scripts/processor.mjs</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/sonseongjun/knowledge-base/logs/processor-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sonseongjun/knowledge-base/logs/processor-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/sonseongjun/.nvm/versions/node/v20.18.1/bin:/Users/sonseongjun/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/sonseongjun</string>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
